# conf.d/websockify.conf
[program:websockify]
command=websockify --web /usr/share/novnc $PORT localhost:5900
autorestart=true
